# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
minfds=1024
minprocs=200

# Give nginx + node a moment to settle before reloader kicks in
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=2
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:node]
; --max-old-space-size=256 keeps Node under Koyeb's 512MB limit
command=node --max-old-space-size=256 /app/server.js
directory=/app
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1",DATABASE_URL="file:/app/db/custom.db"
autostart=true
autorestart=true
startsecs=5
startretries=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; The reloader waits 20s for Next.js to come up, then refreshes the nginx config
; every 5 minutes (instead of every 60s — saves CPU on free tier).
[program:nginx-reloader]
command=sh -c "sleep 20; while true; do curl -sf http://127.0.0.1:3000/api/nginx-config > /etc/nginx/nginx.conf.new && mv /etc/nginx/nginx.conf.new /etc/nginx/nginx.conf && nginx -t 2>/dev/null && nginx -s reload 2>/dev/null; sleep 300; done"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
