# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:node]
command=node /app/server.js
directory=/app
environment=NODE_ENV="production",PORT="3000",HOSTNAME="127.0.0.1",DATABASE_URL="file:/app/db/custom.db"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Reload nginx every 60s to pick up schedule changes
[program:nginx-reloader]
command=sh -c "while true; do sleep 60; curl -s -X POST -H 'Content-Type: application/json' -d '{\"apply\":true}' http://127.0.0.1:3000/api/nginx-config > /dev/null || true; done"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
